# rootfs/etc/s6-overlay/s6-rc.d/71-check-destination-host/up
#!/bin/sh
exec /etc/s6-overlay/scripts/71-check-destination-host
